# systemd/code-server.service
# Reference only. The code-server package ships a template unit `code-server@.service`,
# which 50-code-server.sh enables as `code-server@__DEV_USER__`. Combined with
# `loginctl enable-linger __DEV_USER__`, it runs without an active login session.
# Shown here for reference / hand-management:
[Unit]
Description=code-server (remote-vs-code) for %i
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/usr/bin/code-server
Restart=always
User=%i
Environment=HOME=/home/%i

[Install]
WantedBy=multi-user.target
